# etc/systemd/system/gokill.service
[Unit]
Description=gokill daemon

[Service]
ExecStart=gokill -c /etc/gokill/config.json
Restart=on-failure
Type=simple
